# .apk/upmaru/cineplex/cineplex.initd
#!/sbin/openrc-run

name='cineplex'
description="Compressor Encoder"
supervisor="s6"
s6_service_path="${RC_SVCDIR}/s6-scan/${name}"
command=/var/lib/cineplex/bin/cineplex

depend() {
  need net s6-svscan
}

start_pre() {
  if [ ! -L "${RC_SVCDIR}/s6-scan/${name}" ]; then
    ln -s "/var/lib/${name}/service" "${RC_SVCDIR}/s6-scan/${name}"
  fi
}